# extra/radeon-profile-daemon.service
[Unit]
Description=Radeon-Profile-Qt6 Daemon

[Service]
Type=simple
ExecStart=/usr/bin/radeon-profile-daemon

[Install]
WantedBy=multi-user.target
